# scripts/jlcam-housekeeping.service
[Unit]
Description=JL-CAM Daily Housekeeping

[Service]
Type=oneshot
ExecStart=/usr/local/bin/jlcam_housekeeping.sh
